# sms_gamer_ups_monitor/rootfs/etc/services.d/sms/run
#!/usr/bin/with-contenv bash
exec python3 /app/sms_gamer_ups_monitor.py
